# server/weather/traffic-display-weather.service
[Unit]
Description=Update weather information

[Service]
Type=oneshot
ExecStart=/home/clay/traffic-display-weather/update.py
